# src/ofono.service.in
[Unit]
Requires=dbus.service
After=dbus.service 
Description=Telephony service
After=syslog.target

[Service]
Type=dbus
BusName=org.ofono
ExecStart=@prefix@/sbin/ofonod -n
Restart=always
RestartSec=3
StandardError=null

[Install]
WantedBy=multi-user.target
